# openstack/docker/keystone/supervisord.conf
[supervisord]
nodaemon=true

[program:syncdb]
command=su -s /bin/sh -c "keystone-manage db_sync" keystone
startsecs=0

[program:keystone]
command=/usr/bin/keystone-all
user=keystone
