# openstack/docker/neutron/supervisord.conf
[supervisord]
nodaemon=true

[program:syncdb]
command=su -s /bin/sh -c "nova-manage db sync" nova
startsecs=0

[program:]
command=/usr/bin/neutron-server
user=neutron
