# openstack/docker/nova-controller/supervisord.conf
[supervisord]
nodaemon=true

[program:syncdb]
command=su -s /bin/sh -c "nova-manage db_sync" nova
startsecs=0

[program:nova-controller]
command=/usr/bin/keystone-all
user=keystone
